# DEBIAN/postrm
#!/usr/bin/sh

systemctl disable adblocker
